# cartrasche.service
[Unit]
Description=Cartrsache OnBoot
     
[Service]
ExecStart=/home/jetson/run_cartrasche.sh
     
[Install]
WantedBy=default.target